# other-configs/etc/profile.d/99-coccinelle-path.sh
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/usr/local/lib
export PYTHONPATH=$PYTHONPATH:/usr/local/lib/coccinelle/python
export COCCINELLE_HOME=/usr/local/lib/coccinelle/
