# other-configs/etc/profile.d/99-compiler-path.sh
export PATH="$PATH:/opt/cross-gcc-linux-9/bin/"
