# other-configs/etc/profile.d/99-venv.sh
if [ -f /usr/local/venv/bin/activate ]; then
	source /usr/local/venv/bin/activate
fi
