# proxy-configuration/etc/profile.d/99-corporate-env.sh

export http_proxy=http://webproxy.ext.ti.com:80
export https_proxy=http://webproxy.ext.ti.com:80
export ftp_proxy=http://webproxy.ext.ti.com:80
export no_proxy=.ti.com
export GIT_PROXY_COMMAND=/usr/bin/git-tunnel.sh
